# deploy/plotter.service
[Unit]
Description=Plotter — garden history tracker
After=network.target

[Service]
Type=simple
User=plotter
Group=plotter
WorkingDirectory=/opt/plotter
ExecStart=/opt/plotter/plotter

# Environment — override any of these in a drop-in or EnvironmentFile
Environment=PLOTTER_PORT=8080
Environment=PLOTTER_DB=/opt/plotter/data/plotter.db
Environment=PLOTTER_UPLOAD_DIR=/opt/plotter/data/uploads

# GCS backup — set both to enable cloud backups via backup.sh
# Environment=GOOGLE_APPLICATION_CREDENTIALS=/opt/plotter/gcp-key.json
# Environment=GCS_BUCKET=bucket-name

# Uncomment to load all vars (including secrets) from a file instead:
# EnvironmentFile=/etc/plotter/env

Restart=on-failure
RestartSec=5s

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/plotter/data

[Install]
WantedBy=multi-user.target
